# deploy/gdelt-wildlife.service.example
# Copy to /etc/systemd/system/gdelt-wildlife.service
# Adjust paths, user, and port to match your server (Apache must proxy to the same port).
# For https://SERVERNAME/gdelt-wildlife/ build the SPA with VITE_BASE_PATH=/gdelt-wildlife/
#
#   sudo cp deploy/gdelt-wildlife.service.example /etc/systemd/system/gdelt-wildlife.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now gdelt-wildlife.service

[Unit]
Description=GDELT wildlife map API (Uvicorn / FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/gdelt-wildlife
Environment=PATH=/srv/gdelt-wildlife/.venv/bin
EnvironmentFile=-/srv/gdelt-wildlife/server/.env
ExecStart=/srv/gdelt-wildlife/.venv/bin/uvicorn server.main:app --host 127.0.0.1 --port 8000 --workers 2 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
